# etc/udev/rules.d/99-intel-cpu-watts-fix.rules
SUBSYSTEM=="cpu", ACTION=="add", RUN+="/bin/bash -c 'if grep -q GenuineIntel /proc/cpuinfo; then chmod o+r /sys/class/powercap/intel-rapl:0/energy_uj; fi'"